# cms-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop    cms || true
systemctl disable cms || true
